# roles/cluster/templates/fancontrol-config.service.j2
[Unit]
Description=Generate fancontrol configuration dynamically
Before=fancontrol.service
After=lm-sensors.service

[Service]
Type=oneshot
ExecStart=/usr/local/bin/generate-fancontrol-config.sh
RemainAfterExit=yes

[Install]
WantedBy=multi-user.target
